# infra/systemd/sapphirix-backend.service
[Unit]
Description=Sapphirix backend
After=network.target postgresql.service

[Service]
Type=simple
WorkingDirectory=/var/www/sapphirix/app/current
EnvironmentFile=/var/www/sapphirix/shared/.env
ExecStart=/usr/bin/env pnpm --filter backend start
Restart=always
RestartSec=5
User=www-data
Group=www-data
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target